# files/openmrs-backup-tools/50-kemr.rules
KERNEL=="sd?1", SUBSYSTEM=="block", ATTRS{model}=="USB*", SYMLINK+="backup", RUN+="/opt/openmrs-backup-tools/sync_usb.sh"